# templates/udev/rules.d/99-lowlatency.rules.j2
# Low-latency audio udev rules
# Generated by Ansible for {{ ansible_hostname }}

# Set scheduler for audio devices
KERNEL=="controlC[0-9]*", GROUP="audio", MODE="0664"

# Audio interface optimizations
SUBSYSTEM=="usb", ATTR{idVendor}=="*", ATTR{idProduct}=="*", ATTR{product}=="*Audio*", TAG+="uaccess"

# Set realtime priority for audio group
KERNEL=="rtc0", GROUP="audio", MODE="0664"

# Low-latency kernel optimizations
{% if enable_lowlatency | default(false) %}
# Set CPU governor for audio processes
ACTION=="add", SUBSYSTEM=="cpu", ATTR{cpufreq/scaling_governor}="{{ audio_cpu_governor | default('performance') }}"

# Disable power management for audio devices
ACTION=="add", SUBSYSTEM=="pci", ATTR{class}=="0x040300", ATTR{power/control}="on"

# Set RLIMIT_RTPRIO for audio group
KERNEL=="rtc[0-9]*", GROUP="audio", MODE="0664"
{% endif %}

# USB audio device optimizations
SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="01", ATTR{bInterfaceSubClass}=="02", GROUP="audio"

# MIDI device permissions
KERNEL=="midi[0-9]*", GROUP="audio", MODE="0664"
SUBSYSTEM=="sound", GROUP="audio", MODE="0664"

# Set buffer sizes for low-latency
{% if audio_buffer_size is defined %}
SUBSYSTEM=="sound", ATTR{power/control}="on", RUN+="/bin/echo {{ audio_buffer_size }} > /sys/module/snd/parameters/cards_limit"
{% endif %}