# configs/ansible/roles/power_management/templates/99-power-management.rules.j2
# Power management udev rules
# Generated by Ansible for {{ ansible_hostname }}

# Set power control for PCI devices
ACTION=="add", SUBSYSTEM=="pci", ATTR{power/control}="auto"

# USB device power management
ACTION=="add", SUBSYSTEM=="usb", TEST=="power/control", ATTR{power/control}="auto"
ACTION=="add", SUBSYSTEM=="usb", TEST=="power/autosuspend", ATTR{power/autosuspend}="2"

# Audio device power management
ACTION=="add", SUBSYSTEM=="pci", ATTR{class}=="0x040300", ATTR{power/control}="auto"

# Network device power management
ACTION=="add", SUBSYSTEM=="net", KERNEL=="eth*", RUN+="/usr/bin/ethtool -s %k wol d"
ACTION=="add", SUBSYSTEM=="net", KERNEL=="wlan*", RUN+="/usr/bin/iw dev %k set power_save on"

# SATA link power management
ACTION=="add", SUBSYSTEM=="scsi_host", KERNEL=="host*", ATTR{link_power_management_policy}="med_power_with_dipm"

# Bluetooth device power management
ACTION=="add", SUBSYSTEM=="bluetooth", ATTR{power/control}="auto"

# Graphics device power management
ACTION=="add", KERNEL=="card[0-9]*", SUBSYSTEM=="drm", DRIVERS=="i915", ATTR{power/control}="auto"

# CPU frequency scaling
ACTION=="add", KERNEL=="cpu[0-9]*", SUBSYSTEM=="cpu", ATTR{cpufreq/scaling_governor}="{{ cpu_scaling_governor | default('powersave') }}"

# Backlight power management
ACTION=="add", SUBSYSTEM=="backlight", KERNEL=="*", ATTR{brightness}="{{ backlight_default_brightness | default('15') }}"

# Runtime PM for all devices
ACTION=="add", ATTR{power/control}=="*", ATTR{power/control}="auto"